# rootdir/etc/systemd/system/mce.service.d/10-dsme.conf
[Service]
Environment=DSME_SOCKFILE=/run/dsme.socket